# etc/systemd/system/opoj-empty-ns@.service
[Unit]
Description=Empty namespace holder (runner %I)

[Service]
ExecStart=/bin/sleep infinity
PrivateIPC=true
PrivateNetwork=true
